# playbook/roles/cowrie/templates/cowrie.service.j2
[Unit]
Description=Cowrie SSH/Telnet Honeypot
After=network.target

[Service]
User=cowrie
Group=cowrie
WorkingDirectory=/home/cowrie/cowrie
ExecStart=/usr/bin/authbind --deep /home/cowrie/cowrie/cowrie-env/bin/python3 -m cowrie.core.server
Restart=on-failure

[Install]
WantedBy=multi-user.target
